# deploy/systemd/telegram-codex-bridge.service.example
[Unit]
Description=Telegram Codex Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=codex
WorkingDirectory=/opt/telegram-codex-bridge
EnvironmentFile=/opt/telegram-codex-bridge/.env
ExecStart=/usr/bin/env node /opt/telegram-codex-bridge/server.mjs
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
